# deploy/linode/transfera-client.service
[Unit]
Description=Transfera Next.js static UI
After=network-online.target transfera-api.service
Wants=network-online.target

[Service]
Type=simple
User=root
EnvironmentFile=-/etc/transfera/env
Environment=CLIENT_PORT=3000
# Override in /etc/transfera/env (80 works via python in run-client.sh)
WorkingDirectory=/opt/transfera-client/client
ExecStart=/opt/transfera-client/deploy/linode/run-client.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
